# classes/pebble-project-intro/admin-project/cronjob/project.cron
02 22 24 05 * /net/class/cs50/labs/grading/project/clone.sh may24
02 22 30 05 * /net/class/cs50/labs/grading/project/clone.sh may30
02 22 31 05 * /net/class/cs50/labs/grading/project/clone.sh may31

